# docker/docker/pilot/test/deploy/rest.api.ini
[unix_http_server]
file=/workspace/deploy/supervisord.sock

[supervisord]
logfile=/workspace/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/workspace/deploy/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///workspace/deploy/supervisord.sock

[program:target]
command=/bin/sh -c ./target
directory=/workspace/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/workspace/log/daemon.log
loglevel=info
startsecs=300

